# config/com.worktracker.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.worktracker.frontend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/npm</string>
        <string>run</string>
        <string>dev</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/camdoherty/Library/CloudStorage/OneDrive-InDigitalMarketing/Documents/VSC/work-location-tracker/frontend</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>VITE_API_BASE</key>
        <string>http://localhost:8001</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/camdoherty/Library/CloudStorage/OneDrive-InDigitalMarketing/Documents/VSC/work-location-tracker/logs/frontend.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/camdoherty/Library/CloudStorage/OneDrive-InDigitalMarketing/Documents/VSC/work-location-tracker/logs/frontend.error.log</string>
</dict>
</plist>
